# deploy/agent-team-worker.timer
[Unit]
Description=Agent Team issue consumer timer

[Timer]
OnBootSec=30s
OnUnitActiveSec=10s
Unit=agent-team-worker.service

[Install]
WantedBy=timers.target
